# lib/templates/upstart.conf.erb
start on runlevel [2345]
stop on runlevel [06]

script
<% if @template_arguments[:slave] %>
exec /usr/bin/ttserver -log <%= log_directory %>/ssdata-slave-<%= @template_arguments[:date] %>.log -sid <%= @template_arguments[:sid] %> -thnum <%= @template_arguments[:thnum] %> -mhost <%= @template_arguments[:mhost] %> -mport <%= @template_arguments[:mport] %> -rts <%= log_directory %>/ttserver-<%= @template_arguments[:date] %>.rts <%= data_directory %>/ssdata-slave-<%= @template_arguments[:date] %>.<%= @template_arguments[:format] %>#bnum=<%= @template_arguments[:bnum] %>#apow=<%= @template_arguments[:apow] %>#fpow=<%= @template_arguments[:fpow] %>#ncnum=<%= @template_arguments[:ncnum] %>#xmsiz=<%= @template_arguments[:xmsiz] %>#opts=<%= @template_arguments[:opts] %>
<% else %>
exec /usr/bin/ttserver -log <%= log_directory %>/ssdata-master-<%= @template_arguments[:date] %>.log -sid <%= @template_arguments[:sid] %> -thnum <%= @template_arguments[:thnum] %> <%= data_directory %>/ssdata-master-<%= @template_arguments[:date] %>.<%= @template_arguments[:format] %>#bnum=<%= @template_arguments[:bnum] %>#apow=<%= @template_arguments[:apow] %>#fpow=<%= @template_arguments[:fpow] %>#ncnum=<%= @template_arguments[:ncnum] %>#xmsiz=<%= @template_arguments[:xmsiz] %>#opts=<%= @template_arguments[:opts] %>
<% end %>
end script

respawn

kill timeout 1000000
